# Automatyczny Stock Market 1.0/asebot.service
[Unit]
Description=ASE Trading Platform API
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/asebot-backend/Algorytm Uczenia Kwantowego LLM
EnvironmentFile=/home/admin/asebot-backend/Algorytm Uczenia Kwantowego LLM/.env
ExecStart=/home/admin/asebot-backend/Algorytm Uczenia Kwantowego LLM/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8008 --workers 4 --proxy-headers
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
